# zdsp1d/systemd/zdsp1d.socket
[Socket]
ListenStream=127.0.0.1:6310

[Install]
WantedBy=sockets.target
