# zdsp1d/systemd/zdsp1d.socket
[Socket]
ListenStream=6310

[Install]
WantedBy=sockets.target
