# sec1000d/systemd/sec1000d.socket
[Socket]
ListenStream=6305

[Install]
WantedBy=sockets.target
