# com5003d/systemd/com5003d.socket
[Unit]
ConditionKernelCommandLine=zera_device=com5003

[Socket]
ListenStream=6307

[Install]
WantedBy=sockets.target
