# mt310s2d/systemd/mt581s2d.socket
[Unit]
ConditionKernelCommandLine=zera_device=mt581s2

[Socket]
ListenStream=6307

[Install]
WantedBy=sockets.target
